# etc/services/dcron/run
#!/bin/sh

set -e
. /etc/env

exec crond -f
