# etc/services/postgres/run
#!/bin/sh

set -e
. /etc/env

exec su-exec postgres postgres $POSTGRES_DAEMON_ARGS
